# systemd/system/vyatta-config-default.service.d/migrate.conf
[Service]
ExecStartPre=/opt/vyatta/sbin/vyatta-image-migrate-lb4
